# .config/systemd/user/binance-crypto-position.service
[Unit]
Description=Binance Crypto Alerts - Position Monitor
StartLimitIntervalSec=0

[Service]
Type=simple
WorkingDirectory=/home/moises-beltran-castro/.openclaw/workspace/binance-crypto-alerts/worker
ExecStart=/usr/bin/node /home/moises-beltran-castro/.openclaw/workspace/binance-crypto-alerts/worker/position-monitor.js
Environment=NODE_ENV=production
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
