# .config/systemd/user/binance-crypto-review.timer
[Unit]
Description=Binance Crypto Alerts - Hourly Review

[Timer]
OnBootSec=60s
OnUnitActiveSec=3600s
Unit=binance-crypto-review.service

[Install]
WantedBy=timers.target
